# connectors/codex/deploy/macchiato-codex-connector.service
# Macchiato Codex Connector — systemd user unit (installed by install.sh)
[Unit]
Description=Macchiato Codex Connector
After=network-online.target

[Service]
WorkingDirectory=%h/.macchiato/codex-app
ExecStart=%h/.macchiato/codex-app/node_modules/.bin/tsx src/index.ts
Restart=always
RestartSec=5
Environment=NODE_ENV=production
Environment=PATH=%h/.local/bin:/usr/local/bin:/usr/bin:/bin

[Install]
WantedBy=default.target
